# crontab
* 6-22 * * * cd /var/www && /usr/local/bin/php artisan schedule:run >> /var/log/cron.log 2>&1
* 6-22 * * * date >> /var/log/cron.log
